# com.krypto.tradingbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.krypto.tradingbot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>run_bot_background.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mjcascio/Documents/KryptoTheTradingBot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mjcascio/Documents/KryptoTheTradingBot/launch_agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mjcascio/Documents/KryptoTheTradingBot/launch_agent_error.log</string>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist> 