# com.krypto.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.krypto.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>start_dashboard.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mjcascio/Documents/KryptoTheTradingBot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mjcascio/Documents/KryptoTheTradingBot/logs/dashboard_launch.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mjcascio/Documents/KryptoTheTradingBot/logs/dashboard_launch_error.log</string>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist> 